# microjobboard-api.service
[Unit]
Description=MicroJobBoard .NET API
After=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/var/www/dotnetapi
ExecStart=/usr/bin/dotnet /var/www/dotnetapi/MicroJobBoard.API.dll

# Environment
Environment=ASPNETCORE_ENVIRONMENT=Production
Environment=ASPNETCORE_URLS=http://0.0.0.0:5000
Environment=DOTNET_PRINT_TELEMETRY_MESSAGE=false

# Restart policy
Restart=always
RestartSec=10
KillSignal=SIGINT
SyslogIdentifier=microjobboard-api

# Security
NoNewPrivileges=true
PrivateTmp=true

# Limits
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
